# rootfs.base/etc/inittab
# /etc/inittab
::sysinit:/etc/rcS
::respawn:/bin/sh
::shutdown:/bin/umount -a -r
#::respawn:/bin/sshd -D
#::respawn:/sbin/getty -L ttyS0 115200 vt100
